# ctf/hackthevote2016/200-mount/etc/inittab
# Setup our shell run scripts
::sysinit:/etc/init.d/rcS

# Start an "askfirst" shell on the console (whatever that may be)
::respawn:-/bin/sh

# Stuff to do when restarting the init process
::restart:/sbin/init

# Stuff to do before rebooting
::shutdown:/bin/umount -a -r
#::shutdown:/sbin/swapoff -a

